# etc/systemd/mrhs.service
[Unit]
Description=mrhs

[Service]
Type=simple
ExecStart=/home/pi/mrhs
Restart=always

[Install]
WantedBy = multi-user.target
