# etc/systemd/mrhs.service
[Unit]
Description=mrhs service

[Service]
Type=simple
ExecStart=/home/pi/git/mrhs/mrhsagent
Restart=always

[Install]
WantedBy=multi-user.target
