# etc/systemd/mrhs.service
[Unit]
Description=mrhs service

[Service]
Type=simple
ExecStart=/usr/local/bin/mrhsagent
Restart=always

[Install]
WantedBy=multi-user.target
